# local-runner/linux/systemd/curation-phase2.timer
[Unit]
Description=Curation Phase 2 (JST 04:30)

[Timer]
OnCalendar=*-*-* 04:30:00 Asia/Tokyo
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
