# local-runner/linux/systemd/curation-phase3.timer
[Unit]
Description=Curation Phase 3 (JST 05:30)

[Timer]
OnCalendar=*-*-* 05:30:00 Asia/Tokyo
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
