# local-runner/linux/systemd/curation-phase3.service
[Unit]
Description=Curation Phase 3 (Fact-check + merge)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-%h/.config/curation/runner.env
ExecStart=/bin/bash {{RUN_PHASE}} 3
